# llama-console-gui.service
#
# sudo cp llama-console-gui.service /etc/systemd/system/llama-console-gui.service
# sudo systemctl daemon-reload
# sudo systemctl enable --now llama-console-gui.service
# systemctl status llama-console-gui.service
# journalctl -u llama-console-gui.service -f
#

[Unit]
Description=NiceGUI LLama Console UI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=alvise
Group=alvise
WorkingDirectory=/home/alvise/llama-console
Environment=HOME=/home/alvise
Environment=PATH=/home/alvise/miniconda3/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/alvise/miniconda3/bin/python /home/alvise/llama-console/llama-console-gui.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
